# conf/supervisor.conf
[program:wechat-robot]
command = gunicorn app:app -b localhost:8000
directory = /usr/local/opt/wechat-robot
