# conf/supervisor.conf
[program:api]
command = gunicorn app:app -b localhost:8000
directory = /home/mcxiaoke/apps/api
user = mcxiaoke
